# rootfs/etc/s6-overlay/s6-rc.d/wg-quick/finish
#!/usr/bin/with-contenv sh

# change chdir
cd /config

# run
exec wg-quick down ./$WG_FILE
