# rootfs/etc/s6-overlay/s6-rc.d/wg-quick/run
#!/usr/bin/with-contenv sh

# change chdir
cd /config

# run
exec wg-quick up ./$WG_FILE
